# resources/com.secondbrain.listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.secondbrain.listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/richardyu/Library/Mobile Documents/iCloud~md~obsidian/Documents/Home/01_Projects/Personal/second-brain/archive/legacy-productivity/_scripts/from-para-2026-02-08/start_listener.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/richardyu/Library/Mobile Documents/iCloud~md~obsidian/Documents/Home/01_Projects/Personal/second-brain/archive/legacy-productivity/_scripts/from-para-2026-02-08</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/secondbrain-listener.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/secondbrain-listener-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/richardyu/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
